# autorun.conf
start on startup

exec /opt/droneMesh/autorun.py
